# pkg/etc/supervisor.d/network.ini
[program:hostapd]
command=hostapd /etc/hostapd.moonbow.config
autostart=true
autorestart=true

[program:udhcpd]
command=udhcpd -f -S /etc/udhcpd.moonbow.config
autostart=true
autorestart=true

